# +DummyDevice/Setup/chaprasan.service
[Unit]
Description=Chaprasan IOT Scripts

[Service]
WorkingDirectory=/srv/chaprasan
ExecStart=/srv/chaprasan/startup.sh
User=root
